# systemd/turbooff.service
[Unit]
Description=Turboboost off switch

[Service]
Type=oneshot
ExecStart=/home/azar/.local/bin/turbooff

[Install]
WantedBy=multi-user.target
